# examples/prj/arleux/etc/supervisor/conf.d/rrd_web_service.conf
[program:rrd_web_service]
command=rrd_web_service
user=pi
autostart=true
autorestart=false
environment=PYTHONUNBUFFURED=1
