# examples/prj/arleux/etc/supervisor/conf.d/rrd_recorder.conf
[program:rrd_recorder]
command=rrd_recorder_arleux_mel
user=pi
autostart=true
autorestart=false
environment=PYTHONUNBUFFURED=1
